# collectd/root/etc/services.d/collectd/run
#!/usr/bin/execlineb
/opt/collectd/sbin/collectd -f -C /etc/collectd/collectd.conf
