# riemann/root/etc/services.d/riemann/run
#!/usr/bin/execlineb
/opt/riemann/bin/riemann start
